# cookbooks/supervisor/templates/etc/supervisor/app.conf.erb
[program:<%= @app_name %>]
command=<%= @command %>
directory=<%= @directory %>
autostart=true
autorestart=true
